# base_image/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=nginx

[program:codebox]
command=codebox -p 8001 run /root/exercise
