# base_image/supervisord.conf
[supervisord]
nodaemon=true

[program:tty]
command=tty.js
