# base_image/supervisord.conf
[supervisord]
stopsignal=KILL

[program:hexlet-ide]
command=hexlet-ide -p 8000 -r /usr/src/app 2>&1
