# base_image/supervisord.conf
[supervisord]
nodaemon=true
stopsignal=KILL

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=nginx

[program:hexlet-ide]
command=hexlet-ide -p 8001 -r /root/exercise 2>&1
