# base_image/supervisord.conf
[supervisord]
stopsignal=KILL
loglevel=debug

[program:hexlet-ide]
command=hexlet-ide -p 8000 -r /usr/src/app 2>&1
redirect_stderr=true
